# email-reporting.conf
start on stopped rc RUNLEVEL=[2345]
exec /frontend-email-reporting/frontend-email-reporting-1.0/bin/frontend-email-reporting -mem 512
